# root/etc/services.d/cloudflared/run
#!/bin/sh
# Cloudflared s6 service

if [ -z "$TUNNEL_TOKEN" ]; then
    echo "[WARN] TUNNEL_TOKEN is not set. Starting in DNS Proxy mode for testing/fallback."
    # Listen on port 2000 (metrics/test port) to satisfy CI port check
    exec /usr/local/bin/cloudflared proxy-dns --port 2000 --address 0.0.0.0
fi

echo "[INFO] Starting cloudflared tunnel..."

# Use --no-autoupdate as recommended for containers
# Metrics are bound to 0.0.0.0 for easier health checking
exec /usr/local/bin/cloudflared --no-autoupdate tunnel \
    --metrics ${TUNNEL_METRICS:-0.0.0.0:2000} \
    run --token "$TUNNEL_TOKEN"
